# deploy/eat-it-api.service
[Unit]
Description=Eat it API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/eat-it/app
Environment=NODE_ENV=production
Environment=PORT=3010
Environment=DATABASE_URL=file:../data/eat-it.db
Environment=APP_URL=https://eat-it.space
EnvironmentFile=-/opt/eat-it/app/.env
ExecStartPre=/usr/bin/npm run db:migrate
ExecStart=/usr/bin/node server/index.mjs
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
